# systemd/spark-vllm-standalone.service.template
[Unit]
Description=Spark vLLM Standalone Server (split-node mode, TP=1 per node)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User={{USERNAME}}
WorkingDirectory=/home/{{USERNAME}}

# ── Config load order ─────────────────────────────────────────────
# cluster.env is NOT loaded here — only model + node overrides.
# Load order: model.env → node.env  (later file wins on any key it defines)
# System-level paths are optional (-); user-level paths are authoritative.
EnvironmentFile=-/etc/spark-tools/model.env
EnvironmentFile=/home/{{USERNAME}}/.config/spark-tools/model.env
EnvironmentFile=-/etc/spark-tools/node.env
EnvironmentFile=-/home/{{USERNAME}}/.config/spark-tools/node.env

# ── Pre-start ────────────────────────────────────────────────────
# Remove any stale container from a previous run.
ExecStartPre=-/usr/bin/docker rm -f spark-vllm
# Drop OS page cache to free memory before loading model weights.
# Requires root; the || true makes it non-fatal when running as a normal user.
ExecStartPre=/usr/bin/bash -c 'sync && echo 3 > /proc/sys/vm/drop_caches || true'

# ── Main ─────────────────────────────────────────────────────────
# VLLM_IMAGE and VLLM_SERVE_CMD come from model.env (set by spark-set-model
# via model profiles) and can be overridden per-node in node.env.
#
# VLLM_SERVE_CMD controls the entrypoint style:
#   "vllm serve"  — NGC images (nvidia_entrypoint.sh, need explicit command)
#   ""            — Docker Hub images (vllm serve is the ENTRYPOINT already)
#
# Strip --enforce-eager from VLLM_EXTRA_ARGS at runtime.
# It is only needed for multi-node TP=2 (Triton allocator crash); standalone
# is always TP=1 where CUDA graphs work fine and give ~25-30% better throughput.
ExecStart=/usr/bin/bash -c '\
  VLLM_EXTRA_ARGS="$(echo $VLLM_EXTRA_ARGS | sed "s/--enforce-eager//g")"; \
  docker run \
  --name spark-vllm \
  --network host \
  --gpus all \
  --shm-size 10.24g \
  -v /home/{{USERNAME}}/.cache/huggingface:/root/.cache/huggingface \
  -v /etc/encodings:/etc/encodings \
  -e TIKTOKEN_ENCODINGS_BASE=/etc/encodings \
  -e VLLM_FLASHINFER_MOE_BACKEND=latency \
  ${VLLM_IMAGE} \
  ${VLLM_SERVE_CMD} "$MODEL_NAME" \
    --tensor-parallel-size "$TP_SIZE" \
    --max-model-len "$MAX_MODEL_LEN" \
    --gpu-memory-utilization "$GPU_MEM_UTIL" \
    --port {{PORT}} \
    $VLLM_EXTRA_ARGS \
    --host 0.0.0.0'

ExecStop=/usr/bin/docker stop spark-vllm

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
